# supervisor/conf.d/gunicorn.ini
[program:gunicorn_s]
command=gunicorn -c gunicorn.config.py main:app
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/logs/access.log